# systemd/notistory.service
[Unit]
Description=Notistory notification recorder (D-Bus monitor)
After=graphical-session.target
PartOf=graphical-session.target

# NOTE: This is a template. `install.sh` substitutes __PYTHON__ and __INSTALL_DIR__
# and writes the resolved unit to ~/.config/systemd/user/notistory.service.

[Service]
Type=simple
ExecStart=__PYTHON__ __INSTALL_DIR__/src/notify-logger.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
